# src/etc/rc.d/poudriered.in
#!/bin/sh
#
# PROVIDE: poudriered
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
#
# poudriered_enable="YES"
#

. /etc/rc.subr

name=poudriered
rcvar=poudriered_enable

# read configuration
load_rc_config ${name}

: ${poudriered_enable:=NO}
: ${poudriered_pidfile:=/var/run/${name}.pid}
: ${poudriered_socket:=/var/run/${name}.sock}
: ${poudriered_conf:=@prefix@/etc/${name}.conf}

pidfile="${poudriered_pidfile}"
command=@prefix@/libexec/poudriere/${name}
extra_commands="reload"

# run service
run_rc_command "$1"
