# src/etc/rc.d/poudriere.in
#!/bin/sh
# PROVIDE: poudriere
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
#
# poudriere_enable="YES"
#

. /etc/rc.subr

name=poudriere
rcvar=poudriere_enable

# read configuration
load_rc_config ${name}

: ${poudriere_enable:=NO}
: ${poudriere_pidfile:=/var/run/${name}.pid}

pidfile="${poudriere_pidfile}"
command=/usr/local/share/${name}/daemon.sh
command_interpreter=/bin/sh
command_args="daemon -p ${pidfile}"

# run service
run_rc_command "$1"
